# bin/recorder.service
[Unit]
Description=Auto Farm Recorder
After=network.target

[Service]
ExecStart=/usr/local/opt/python-3.8.0/bin/python3.8 -u /home/pi/auto_farm_recorder/src/main.py
WorkingDirectory=/home/pi
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=auto_farm_recorder
Restart=always
RestartSec=1
User=root

[Install]
WantedBy=multi-user.target